# roles/database/templates/backup-database.cron.j2
0 1 * * * python3 {{ script_root }}/{{ project_name }}/backup-database.py