# roles/database/templates/backup-database.cron.j2
0 1 * * * root cd {{ script_root }}/{{ project_name }} && python3 backup-database.py >> /var/log/backend/backup-database.log 2>&1